# 202.38.96.201/rc.local
#!/bin/sh

#centos 5

echo 1 > /proc/sys/net/ipv4/ip_forward

modprobe ip_conntrack hashsize=200000
modprobe ip_conntrack_ftp
modprobe ipt_state
modprobe iptable_nat
modprobe ip_nat_ftp

cd /proc/sys/net/ipv4/netfilter

echo 60 > ip_conntrack_generic_timeout
echo 10 > ip_conntrack_icmp_timeout
echo 10 > ip_conntrack_tcp_timeout_close_wait
echo 900 > ip_conntrack_tcp_timeout_established
echo 10 > ip_conntrack_tcp_timeout_fin_wait
echo 10 > ip_conntrack_tcp_timeout_last_ack
echo 10 > ip_conntrack_tcp_timeout_syn_recv
echo 10 > ip_conntrack_tcp_timeout_syn_sent
echo 10 > ip_conntrack_tcp_timeout_time_wait



ip add add 218.104.71.162/28 dev eth0

ip route add 202.38.64.0/19 via 202.38.96.194
ip route add 210.45.64.0/20 via 202.38.96.194
ip route add 210.45.112.0/20 via 202.38.96.194
ip route add 211.86.144.0/20 via 202.38.96.194
ip route add 222.195.64.0/19 via 202.38.96.194
ip route add 114.214.160.0/19 via 202.38.96.194
ip route add 114.214.192.0/19 via 202.38.96.194


/etc/rc.d/rc.firewall
